# odinforge-agent/deploy/systemd/odinforge-agent.service
[Unit]
Description=OdinForge Security Agent
Documentation=https://github.com/odinforge/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=odinforge
Group=odinforge

ExecStart=/usr/local/bin/odinforge-agent --config /etc/odinforge/agent.yaml
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

# Environment file for secrets
EnvironmentFile=-/etc/odinforge/agent.env

# State and runtime directories
StateDirectory=odinforge-agent
RuntimeDirectory=odinforge-agent
ConfigurationDirectory=odinforge

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/odinforge-agent

# Limit capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native

# Resource limits
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
