# odinforge-agent/deploy/launchd/com.odinforge.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.odinforge.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/odinforge-agent</string>
        <string>--config</string>
        <string>/etc/odinforge/agent.yaml</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>ODINFORGE_SERVER_URL</key>
        <string>https://your-odinforge-server.com</string>
        <key>ODINFORGE_API_KEY</key>
        <string>your-api-key-here</string>
        <key>ODINFORGE_TENANT_ID</key>
        <string>default</string>
        <key>ODINFORGE_AUTH_MODE</key>
        <string>api_key</string>
        <key>ODINFORGE_REQUIRE_HTTPS</key>
        <string>true</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>WorkingDirectory</key>
    <string>/var/lib/odinforge-agent</string>
    
    <key>StandardOutPath</key>
    <string>/var/log/odinforge-agent/agent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/odinforge-agent/agent.error.log</string>
    
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
